# amdgpu-fan.service
[Unit]
Description=amdgpu fan controller

[Service]
ExecStart=/usr/bin/amdgpu-fan
Restart=always

[Install]
WantedBy=sysinit.target
